# config/crontab
*/10 * * * * cd /myapp && bundle exec rails build:all && bundle exec rails build:rsync
